# templates/default/varnish.conf.erb
description "varnish"

start on net-device-up
stop on runlevel [!2345]

respawn
respawn limit 5 30

limit nofile <%= node[:varnishd][:limit][:nofile] %> <%= node[:varnishd][:limit][:nofile] %>   
limit memlock <%= node[:varnishd][:limit][:memlock] %> <%= node[:varnishd][:limit][:memlock] %> 

exec /usr/local/sbin/varnishd\
  -F\
  -S /usr/local/etc/varnish/secret\
  -s <%= node[:varnishd][:storage] %>\
  -f /usr/local/etc/varnish/default.vcl\
  -P <%= node[:varnishd][:pid] %>\
<% node[:varnishd][:runtime].each_pair do |key, value| %>
<% unless value.nil? %>
  -p <%= key %>=<%= value %>\
<% end %>
<% end %>
